# network-tools/wint-off@.service
[Unit]
Description=Turn off wint after %i minutes

[Service]
Type=oneshot
ExecStartPre=/bin/sleep %i
ExecStart=/bin/ip link set wint down

